# shared/vlfdaemon.conf
#
# Upstart script for ina-server
#

description "The VLF observatory network daemon"

exec /usr/bin/vlfdaemon

start on runlevel [23]
stop on runlevel [016]

